# scripts/systemd/way-esports.service
[Unit]
Description=WAY Esports stack (docker compose)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/WAY-Esports
ExecStartPre=/usr/bin/git pull --rebase --autostash || true
ExecStart=/usr/bin/docker compose up -d --build
ExecReload=/usr/bin/docker compose up -d --build
ExecStop=/usr/bin/docker compose down
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target


